# ansible/roles/application/templates/.service
[Unit]
Description=Gunicorn instance to serve edu-flask-app
After=network.target

[Service]
User=ubuntu
Group=www-data
WorkingDirectory=/home/ubuntu/edu-flask-app
Environment="PATH=/home/ubuntu/edu-flask-app/venv/bin"
Environment="database_user={{ database_user }}"
Environment="database_password={{ database_password }}"
Environment="database_name={{ database_name }}"
Environment="database_host={{ database_host }}"
ExecStart=/home/ubuntu/edu-flask-app/venv/bin/gunicorn --workers 3 --bind 0.0.0.0:5000 wsgi:app

[Install]
WantedBy=multi-user.target
