# ansible/roles/application/templates/.service
[Unit]
Description=Gunicorn instance to serve {{ app_name }}
After=network.target

[Service]
User=ubuntu
Group=www-data
WorkingDirectory=/home/ubuntu/{{ app_name }}
Environment="PATH=/home/ubuntu/{{ app_name }}/venv/bin"
ExecStart=/home/ubuntu/{{ app_name }}/venv/bin/gunicorn main:app --workers 3 --worker-class uvicorn.workers.UvicornWorker --bind 0.0.0.0:8000 

[Install]
WantedBy=multi-user.target
